# conf/supervisor/icecast.ini
[program:icecast]
command=/usr/bin/icecast -c /config/icecast.xml
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true